# tools/openmoss/ops/ai.jinclaw.neosgo-lead-collection-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.jinclaw.neosgo-lead-collection-daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>/Users/mac_claw/.openclaw/workspace/projects/neosgo-marketing-suite/scripts/send_lead_collection_progress_telegram.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/mac_claw/.openclaw/workspace/output/neosgo-outreach/lead-collection.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mac_claw/.openclaw/workspace/output/neosgo-outreach/lead-collection.stderr.log</string>
</dict>
</plist>
